# vm/rootfs/etc/systemd/system/destroy-vm.service
[Unit]
Description="Systemd service that invokes the script that destroys a VM"

[Service]
ExecStart=/etc/systemd/system/shutdown.sh
